# packaging/launchpad-process-pool.socket
[Socket]
ListenStream=/run/user/%U/.launchpad-process-pool-sock

Service=launchpad-process-pool.service

[Install]
WantedBy=sockets.target
